# zabbix-agent/netstat/etc/cron.d/sysstat
# statystyka netstat do cache
* * * * * root /usr/sbin/netstat-cache-fetch

